# .docker/cron/crontab
0 * * * * /usr/local/bin/php /var/www/crypto/bin/console run:sync-rates >> /var/log/cron.log 2>&1
